# systemd/.config/systemd/user/wireplumber.service.d/override.conf
[Service]
LimitRTPRIO=95
LimitNICE=-20
LimitMEMLOCK=8388608
Nice=-15
NoNewPrivileges=no
