# systemd/.config/systemd/user/pipewire.service.d/override.conf
[Service]
OOMScoreAdjust=-100
LimitNICE=-20
LimitRTPRIO=95
LimitMEMLOCK=8388608
Nice=-15
Environment=PIPEWIRE_RTKIT=0
