# systemd/.config/systemd/user/pipewire-pulse.service.d/override.conf
[Service]
LimitNICE=-11
LimitRTPRIO=95
LimitMEMLOCK=8388608
