# systemd/.config/systemd/user/bing-wallpaper.timer
[Unit]
Description=Download Bing wallpaper daily

[Timer]
OnCalendar=daily
OnBootSec=5min
Persistent=true

[Install]
WantedBy=timers.target
