# systemd/.config/systemd/user/bluetooth-autoconnect.service
[Unit]
Description=Auto-connect trusted Bluetooth audio devices
After=pipewire.service wireplumber.service bluetooth.service
Requires=pipewire.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/bash -c 'for i in $(seq 1 10); do bluetoothctl connect 5C:02:14:2C:E6:D6 && exit 0; sleep 3; done; exit 1'

[Install]
WantedBy=default.target
