# systemd/.config/systemd/user/hermes-gateway.service
[Unit]
Description=Hermes Agent Gateway - Messaging Platform Integration
After=network.target
StartLimitIntervalSec=600
StartLimitBurst=5

[Service]
Type=simple
ExecStart=/home/wert/.hermes/hermes-agent/venv/bin/python -m hermes_cli.main gateway run --replace
WorkingDirectory=/home/wert/.hermes/hermes-agent
Environment="PATH=/home/wert/.hermes/hermes-agent/venv/bin:/home/wert/.hermes/hermes-agent/node_modules/.bin:/home/wert/.nvm/versions/node/v25.8.1/bin:/home/wert/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="VIRTUAL_ENV=/home/wert/.hermes/hermes-agent/venv"
Environment="HERMES_HOME=/home/wert/.hermes"
Restart=on-failure
RestartSec=30
RestartForceExitStatus=75
KillMode=mixed
KillSignal=SIGTERM
ExecReload=/bin/kill -USR1 $MAINPID
TimeoutStopSec=90
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
