# systemd/.config/systemd/user/mpd-crossfade.service
[Unit]
Description=MPD Crossfade Configuration
After=mpd.service
Wants=mpd.service

[Service]
Type=oneshot
ExecStart=%h/.config/mpd/set-crossfade.sh
RemainAfterExit=yes

[Install]
WantedBy=default.target