# systemd/.config/systemd/user/wallpaper-rotate.service
[Unit]
Description=Rotate local wallpaper randomly

[Service]
Type=oneshot
ExecStart=%h/.local/scripts/wallpaper-rotate.sh
